# dist/hptdrv-monitor
#!/bin/bash

# chkconfig: 12345 10 10
# description: Rebuild driver for HighPoint HBA if kernel changed.
### BEGIN INIT INFO
# Provides: hptdrv-monitor
# Required-Start:
# Required-Stop: reboot network networking network-manager umountroot
# Default-Start: 1 2 3 4 5
# Default-Stop: 0 6
# Short-Description: Rebuild driver for HighPoint HBA if kernel changed.
### END INIT INFO

[ -f /etc/init.d/functions ] && . /etc/init.d/functions

quiet=1
buildall=1

DAEMON=hptdrv-monitor
pidfile=/var/run/yum.pid
drvroot=/usr/share/hptdrv
if test ! -d "${drvroot}"; then
  exit 0
fi

goon()
{
  logecho "I am trapped"
  continue
}

trap goon SIGHUP SIGINT SIGTERM

. ${drvroot}/hptdrv-function

if test "${modules}" = ""; then
  exit 0
fi

temp=/var/tmp
host=www.highpoint-tech.com
webroot=http://${host}/linux

logfile=/var/log/hptdrv.log
lock=/var/lib/hptdrv-monitor
[ -d /var/lock/subsys ] && lock=/var/lock/subsys/hptdrv-monitor

mkinitrd=`type -p dracut`
[ "$mkinitrd" = "" ] && mkinitrd=`type -p mkinitramfs`
[ "$mkinitrd" = "" ] && mkinitrd=`type -p mkinitrd`

rtrn=1
warn=1
FORCE=0
autoup=1

test "$runlevel" = "0" && warn=0
test "$runlevel" = "6" && warn=0

if test "${kernels}" = ""; then
  logecho "no kernels detected."
  exit 0
fi

mkdir -p "${temp}"

case "$1" in
  start)
    touch ${lock}
    start
    rtrn=$?
    ;;
  stop)
    if test -f /var/lib/hptdrv.rebuild; then
      FORCE=1
      export FORCE
    fi
    rm -f "${logfile}"
    compile
    rtrn=$?
    if test "$rtrn" != "0"; then
      update
      rtrn=$?
      if test "$rtrn" = "0"; then
        compile
        rtrn=$?
      fi
    fi
    fixboot
    if test -f /var/lib/hptdrv.rebuild; then
      rm -f /var/lib/hptdrv.rebuild
    fi
    ;;
  restart)
    touch ${lock}
    rtrn=0
    ;;
  status)
    rtrn=0
    ;;
  *)
    echo $"Usage: $0 {start|stop|restart|status}"
    ;;
esac

exit $rtrn

# vim: expandtab ts=2 sw=2 ai
